# infra/hetzner/ops/systemd/pluginscore-docker-prune.service
[Unit]
Description=Prune PluginScore Docker build cache
Wants=docker.service
After=docker.service

[Service]
Type=oneshot
ExecStart=/opt/pluginscore/infra/hetzner/ops/pluginscore-docker-prune.sh
