# infra/hetzner/ops/systemd/pluginscore-docker-prune.timer
[Unit]
Description=Run PluginScore Docker build-cache prune

[Timer]
OnCalendar=Sun *-*-* 03:30:00
Persistent=true
RandomizedDelaySec=20min

[Install]
WantedBy=timers.target
